# backend/agents/Features/Deploy_management/Templates/civilai-agents.service
[Unit]
Description=CivilAI Agentic Engineering Checks
Documentation=file:/home/YOUR_LINUX_USER/CivilAI/backend/agents/README.md
Wants=network-online.target
After=network-online.target civilai-backend.service

[Service]
Type=oneshot
User=YOUR_LINUX_USER
WorkingDirectory=/home/YOUR_LINUX_USER/CivilAI
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/etc/civilai/agents.env
ExecStart=/home/YOUR_LINUX_USER/CivilAI/backend/agents/run_agents_server.sh

# Keep the report run bounded so a stuck network dependency cannot pile up.
TimeoutStartSec=20min

[Install]
WantedBy=multi-user.target
