# deploy/civilai-backend.service
[Unit]
Description=Civil AI Backend
After=network.target

[Service]
User=YOUR_LINUX_USER
WorkingDirectory=/home/YOUR_LINUX_USER/CivilAI
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/YOUR_LINUX_USER/CivilAI/.venv/bin/uvicorn backend.Features.Runtime_management.backend_run:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
